# install/systemd/pi-rtk-surveyor.service
[Unit]
Description=Pi RTK Surveyor Application
Documentation=https://github.com/hotschmoe/pi-rtk-surveyor
After=network.target
Wants=network.target

[Service]
Type=simple
User=PLACEHOLDER_USER
Group=PLACEHOLDER_USER
WorkingDirectory=PLACEHOLDER_PROJECT_DIR/software
ExecStart=PLACEHOLDER_PROJECT_DIR/venv/bin/python PLACEHOLDER_PROJECT_DIR/software/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=PLACEHOLDER_PROJECT_DIR/software
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=PLACEHOLDER_PROJECT_DIR/data

# Resource limits
MemoryLimit=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
