# install/systemd/pi-rtk-surveyor.service
[Unit]
Description=Pi RTK Surveyor Application
Documentation=https://github.com/hotschmoe/pi-rtk-surveyor
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-rtk-surveyor/software
ExecStart=/usr/bin/python3 /home/pi/pi-rtk-surveyor/software/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/pi-rtk-surveyor/software
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/pi-rtk-surveyor/data

# Resource limits
MemoryLimit=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
